# mcp/deploy/loopclub-mcp-http.service
# loopclub-mcp remote (Streamable HTTP) — systemd USER service.
#
# Matches the existing tg-itsavibe.com service pattern on the VPS (nvm node,
# user units, linger enabled). No Docker, no sudo for the service itself.
#
#   cp deploy/loopclub-mcp-http.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now loopclub-mcp-http
#   systemctl --user status loopclub-mcp-http
#
# (Linger is already on for theo, so the service survives logout / reboot.)

[Unit]
Description=loopclub MCP (remote Streamable HTTP)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/theo/projects/loopclub/mcp
Environment="PATH=/home/theo/.nvm/versions/node/v22.22.0/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/theo/.nvm/versions/node/v22.22.0/bin/node dist/http.js
Restart=always
RestartSec=5

# ── Configuration ──────────────────────────────────────────────────────────
# Port 3008 — next free slot after the existing 3001–3007 services.
Environment=PORT=3008
Environment=MCP_BIND_HOST=127.0.0.1
# Accept the tunnel hostname (mcp.tg-itsavibe.com) plus the branded host if/when
# loopclub.xyz DNS moves to Cloudflare, plus localhost for health checks.
Environment=MCP_ALLOWED_HOSTS=mcp.tg-itsavibe.com,mcp.loopclub.xyz,localhost,127.0.0.1
Environment=MCP_ALLOWED_ORIGINS=https://app.loopclub.xyz,https://loopclub.xyz
Environment=LOOPCLUB_ORIGIN=https://app.loopclub.xyz
# Optional tuning (defaults shown):
# Environment=MCP_MAX_BODY_BYTES=65536
# Environment=MCP_RATE_MAX=120
# Environment=MCP_RATE_WINDOW_MS=60000

# ── Containment (the noisy-neighbor blast-radius fence) ──────────────────────
# Caps so a flood that slips past Cloudflare can't starve the co-located
# services (telegram-agent, shillscore, vitals, …). User-unit resource control
# works on Ubuntu 24.04 with the default cgroup delegation; if a controller
# isn't delegated, systemd ignores the directive (non-fatal).
NoNewPrivileges=true
PrivateTmp=true
MemoryMax=256M
TasksMax=64
CPUQuota=80%

[Install]
WantedBy=default.target
